# CredentialHarvestingDetector/credential-harvesting-detector.service
[Unit]
Description=Credential Harvesting Detector Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/credential_harvesting_detector.py
Restart=on-failure
RestartSec=60s

# Resource limits for desktop usage
CPUQuota=10%
CPUWeight=50
MemoryMax=150M
IOWeight=50
TasksMax=64

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
MemoryDenyWriteExecute=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes
SystemCallFilter=~@clock @debug @cpu-emulation @keyring @module @mount @obsolete @raw-io

[Install]
WantedBy=multi-user.target
