# AnomalyCPUMemHunter/anomaly-hunter.service
[Unit]
Description=Anomaly CPU/Memory Usage Hunter Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/anomaly_hunter.py /etc/anomaly_hunter.conf
Restart=on-failure
RestartSec=10
User=root
Group=root
# Security hardening
ProtectSystem=full
ProtectHome=read-only
PrivateDevices=true
NoNewPrivileges=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
