# ZombieHunter/zombie-hunter.service
[Unit]
Description=Zombie Hunter - Suspicious Process Watchdog
After=network.target
Documentation=https://securonis.local/docs/zombie_hunter

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/zombie_hunter.py
Restart=on-failure
RestartSec=30s
KillMode=process
KillSignal=SIGTERM

# Security hardening for the monitor itself
ProtectSystem=strict
ReadWritePaths=/var/log/zombie_hunter /var/run
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true

# Needs to be able to see all processes to work correctly
# But restrict capabilities to only what's needed
CapabilityBoundingSet=CAP_KILL CAP_DAC_READ_SEARCH CAP_SYS_PTRACE

[Install]
WantedBy=multi-user.target
